# ci/build/mcode-server-user.service
[Unit]
Description=mcode-server
After=network.target

[Service]
Type=exec
ExecStart=/usr/bin/mcode-server
Restart=always

[Install]
WantedBy=default.target
